# infrastructure/jobs/crontab
*/1 * * * * /opt/jobs/scripts/github-gitlab-sync.sh >> /var/log/github-gitlab-sync.log 2>&1
