# infrastructure/jobs/crontab
*/1 * * * * /opt/jobs/cron-configs/github-gitlab-sync/script.sh >> /var/log/github-gitlab-sync.log 2>&1
